# crates/velnor-runner/debian/velnor-doctor@.service
[Unit]
Description=Velnor fleet health probe (%i)
Documentation=https://github.com/tailrocks/velnor
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/velnor/%i.env
EnvironmentFile=-/etc/velnor/%i.secrets.env
ExecStart=/usr/bin/velnor-runner doctor \
  --url ${VELNOR_URL} \
  --name ${VELNOR_NAME} \
  --slots ${VELNOR_SLOTS}
